# ansible/inventory/resources/backup/gitlab-pod-restore.service
[Unit]
After=default.target
Requires=default.target

[Service]
ExecStart=/usr/local/bin/gitlab-pod-restore.sh
Type=simple
